# i18n/korean/root/etc/uci-defaults/luci-i18n-korean
#!/bin/sh
uci batch <<-EOF
	set luci.languages.ko=한국어
​        commit luci
EOF

